# tcplog_dumper.in
#!/bin/sh

# PROVIDE: tcplog_dumper
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=tcplog_dumper
rcvar=tcplog_dumper_enable

load_rc_config ${name}

: ${tcplog_dumper_enable:=NO}
: ${tcplog_dumper_pidfile=/var/run/${name}.pid}
: ${tcplog_dumper_basedir=/var/log/tcplog_dumps}
: ${tcplog_dumper_basedir_owner=nobody}

command="/usr/local/sbin/${name}" 
command_args="-d -D ${tcplog_dumper_basedir} -u ${tcplog_dumper_basedir_owner} -p ${tcplog_dumper_pidfile}" 

pidfile="${tcplog_dumper_pidfile}"
start_precmd="${name}_prestart"

tcplog_dumper_prestart()
{
	if [ ! -d "${tcplog_dumper_basedir}" ]; then
		mkdir -p "${tcplog_dumper_basedir}"
	fi
	if [ ! "x$(stat -f '%Su' ${tcplog_dumper_basedir})" = "x${tcplog_dumper_basedir_owner}" ] ; then
		chown "${tcplog_dumper_basedir_owner}" "${tcplog_dumper_basedir}"
	fi
}

run_rc_command "$1"
